# etc/supervisor/conf.d/mangosd.conf
[program:mangosd]
command=/home/mangos/mangos/run/bin/mangosd
directory=/home/mangos/mangos/run/bin/
user=mangos
priority=999
