# etc/supervisor/conf.d/realmd.conf
[program:realmd]
command=/home/mangos/mangos/run/bin/realmd
directory=/home/mangos/mangos/run/bin/
user=mangos
priority=1
