# etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd
directory=/usr/sbin/
user=root
priority=10
